# rngd.service.in
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon
Documentation=man:rngd(8)

[Service]
ExecStart=@sbindir@/rngd -f

[Install]
WantedBy=multi-user.target
